# iso-profiles/shared/garuda/live-overlay-common/etc/xdg/autostart/garuda-system-maintenance.desktop
[Desktop Entry]
Exec=/usr/bin/false
Icon=garuda-system-maintenance
Name=Garuda System Maintenance
Hidden=true
Type=Application
